# SyslogFlipswitch/syslogflipswitchserver/layout/DEBIAN/postinst
#!/bin/sh
launchctl load /Library/LaunchDaemons/de.j-gessner.syslogflipswitchxpc.plist
